# deploy/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

# ============================================================
# NexusBox Sandbox - Supervisor Process Manager
# ============================================================
# Service startup order (by priority):
#   1. entrypoint  - Environment validation & dependency checks
#   2. vnc         - Virtual framebuffer (Xvnc)
#   3. nexusbox    - Gateway + MCP Hub + Port Proxy (single process)
#   4. websockify  - VNC WebSocket bridge (port 6080)
#   5. jupyter     - JupyterLab (port 8888)
#   6. code-server - VS Code Server (port 8200)
#   7. chromium    - Headless browser with CDP (port 9222)
# ============================================================

# --- Entrypoint: environment validation ---
[program:entrypoint]
command=/usr/local/bin/nexusbox-entrypoint.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=1

# --- VNC: Virtual framebuffer ---
[program:vnc]
command=/usr/bin/Xtigervnc :0 -geometry 1280x1024 -depth 24 -SecurityTypes None
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/vnc-stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/vnc-stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=5

# --- NexusBox: Gateway + MCP Hub + Port Proxy ---
# All three services run in a single process (sandbox-dev)
[program:nexusbox]
command=/usr/local/bin/nexusbox-agent -port=8080 -mcp-port=8079 -proxy-port=6081 -workspace=/home/sandbox -log-level=%(ENV_LOG_LEVEL)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/nexusbox-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/nexusbox-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=10

# --- WebSocket Proxy: VNC bridge ---
[program:websockify]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/websockify-stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/websockify-stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=15

# --- JupyterLab: Python notebook ---
[program:jupyter]
command=/usr/local/bin/jupyter lab --no-browser --ip=0.0.0.0 --port=8888 --NotebookApp.token='' --NotebookApp.password='' --NotebookApp.allow_origin='*'
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/jupyter-stdout.log
stdout_logfile_maxbytes=30MB
stdout_logfile_backups=4
stderr_logfile=/var/log/supervisor/jupyter-stderr.log
stderr_logfile_maxbytes=30MB
stderr_logfile_backups=4
priority=25

# --- Code Server: VS Code in browser ---
[program:code-server]
command=/usr/local/bin/code-server --bind-addr 0.0.0.0:8200 --auth none /home/sandbox
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/code-server-stdout.log
stdout_logfile_maxbytes=30MB
stdout_logfile_backups=4
stderr_logfile=/var/log/supervisor/code-server-stderr.log
stderr_logfile_maxbytes=30MB
stderr_logfile_backups=4
priority=25

# --- Chromium: Headless browser with CDP ---
[program:chromium]
command=/usr/bin/chromium --no-sandbox --disable-gpu --disable-dev-shm-usage --remote-debugging-port=9222 --remote-debugging-address=0.0.0.0 --window-size=1280x1024 --user-data-dir=/tmp/chromium http://localhost:8080
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/chromium-stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/chromium-stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=30

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
